# queueapp/queue.conf
[program:queue]
command=php artisan queue:work
directory=/var/www/html
redirect_stderr=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
